# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:redis]
command=redis-server --save "" --appendonly no
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true

[program:auth-server]
command=/usr/local/bin/auth-server
autostart=true
autorestart=true
environment=DATABASE_URL="%(ENV_DATABASE_URL)s",REDIS_URL="redis://127.0.0.1:6379",JWT_SECRET="%(ENV_JWT_SECRET)s"

[program:project-manager]
command=/usr/local/bin/project-manager
autostart=true
autorestart=true
environment=DATABASE_URL="%(ENV_DATABASE_URL)s",REDIS_URL="redis://127.0.0.1:6379",JWT_SECRET="%(ENV_JWT_SECRET)s"

[program:db-proxy]
command=/usr/local/bin/db-proxy
autostart=true
autorestart=true
environment=PG_TARGET="%(ENV_PG_TARGET)s"

[program:storage-api]
command=/usr/local/bin/storage
autostart=true
autorestart=true
environment=DATABASE_URL="%(ENV_DATABASE_URL)s"

[program:edge-functions]
command=/usr/local/bin/edge-functions
autostart=true
autorestart=true

[program:sql-editor]
command=/usr/local/bin/sql-editor
autostart=true
autorestart=true
environment=DATABASE_URL="%(ENV_DATABASE_URL)s"

[program:ai-assistant]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 3006
directory=/app/ai-assistant
autostart=true
autorestart=true
environment=MISTRAL_API_KEY="%(ENV_MISTRAL_API_KEY)s"
