# supervisord.conf
[supervisord]
nodaemon=true


[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=2

[program:redis]
command=redis-server --save "" --appendonly no
autostart=true
autorestart=true
priority=2

[program:auth]
command=/usr/local/bin/auth-server
autostart=true
autorestart=true
priority=3

[program:project-manager]
command=/usr/local/bin/project-manager
autostart=true
autorestart=true
priority=3

[program:db-proxy]
command=/usr/local/bin/db-proxy
autostart=true
autorestart=true
priority=3

[program:storage]
command=/usr/local/bin/storage
autostart=true
autorestart=true
priority=3

[program:sql-editor]
command=/usr/local/bin/sql-editor
autostart=true
autorestart=true
priority=3

[program:edge-functions]
command=/usr/local/bin/edge-functions
autostart=true
autorestart=true
priority=3

[program:ai-assistant]
command=python3 /app/ai-assistant/server.py
autostart=true
autorestart=true
priority=3

[program:db-init]
command=bash -c "while ! pg_isready -q -h 127.0.0.1 -p 5432; do sleep 1; done; if ! psql -h 127.0.0.1 -U blubase -lqt | cut -d \| -f 1 | grep -qw blubase_control; then createdb -h 127.0.0.1 -U blubase blubase_control && psql -h 127.0.0.1 -U blubase blubase_control -f /app/seed.sql; fi"
autostart=true
autorestart=false
startsecs=0
